# Day3/pluginconf-srv/systemd/tomcat-webgw.service
[Unit]
Description=Tomcat web-tier gateway (plugin-conf lab, HTTP 8080)
After=network.target

[Service]
Type=forking
User=tomcat
Group=tomcat
Environment=JAVA_HOME=/usr/lib/jvm/java-21-openjdk-amd64
Environment=CATALINA_HOME=/opt/tomcat11
Environment=CATALINA_BASE=/srv/webgw
Environment=CATALINA_PID=/srv/webgw/temp/tomcat.pid
ExecStart=/opt/tomcat11/bin/catalina.sh start
ExecStop=/opt/tomcat11/bin/catalina.sh stop

[Install]
WantedBy=multi-user.target
